# rootfs/etc/services.d/dab_zrx_watcher-sidecar/run
#!/usr/bin/with-contenv bash
set -e

. /etc/common.sh

echo "Waiting for consul..."
s6-svwait -U /var/run/s6/services/consul

exec consul connect proxy -sidecar-for dab_zrx_watcher
